# docker/cron/etc/cron.d/cron
* * * * * echo $(date +\%Y\-\%m\-\%d\ \%H\:\%M) "- Users checked." >> /var/log/cron.log
# Empty line